# jv_restart_job.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/jv_restart_job.sh

[Install]
WantedBy=default.target
